# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 --nodaemonize
autostart=true
autorestart=true
redirect_stderr=true

[program:mysqld]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
redirect_stderr=true
