# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:mysql]
command=/usr/bin/mysqld_safe

[program:php-fpm]
command=/usr/sbin/php-fpm7.2 --nodaemonize
